# manifests/50-vagrant-mount.rules
SUBSYSTEM=="bdi",ACTION=="add",RUN+="/etc/init.d/supervisord start"
SUBSYSTEM=="bdi",ACTION=="remove",RUN+="/etc/init.d/supervisord stop"
